# mps-device-emulator-rpm/src/main/scripts/postinstall
croot=/usr/share/mps-device-emulator

# logs to /var/log
mkdir -p /var/log/mps-device-emulator
chmod +w /var/log/mps-device-emulator

chmod +x $croot/bin/start
chmod +x $croot/bin/stop

# add vsys init scripts for Ubuntu/CentOS
if [ -x /usr/sbin/update-rc.d ]; then
        update-rc.d mps-device-emulator defaults 86 14
elif [ -x /sbin/chkconfig ]; then
        chkconfig --add mps-device-emulator
fi

service mps-device-emulator start 

