# mps-device-emulator-rpm/src/main/scripts/preremove
# stop service
service mps-device-emulator stop
sleep 5

# clean up folders
rm -rf /var/log/mps-device-emulator
rm -rf /usr/share/mps-device-emulator

# remove vsys init scripts for Ubuntu/CentOS
if [ -x /usr/sbin/update-rc.d ]; then
    update-rc.d -f mps-device-emulator remove
elif [ -x /sbin/chkconfig ]; then
    chkconfig --del mps-device-emulator
fi
