# linux-agent/assembly/src/main/base-deb/control/postrm
if [ "$1" = "purge" ] ; then
    update-rc.d cumulocity-agent remove > /dev/null
fi
