# linux-agent/assembly/src/main/base-deb/control/prerm
if [ -x "/etc/cumulocity-agent.properties" ]; then
    rm /etc/cumulocity-agent.properties
fi

if [ -x "/etc/init.d/cumulocity-agent" ]; then
    invoke-rc.d cumulocity-agent stop || true
fi
